# deploy/goes-nc-cache-cleanup.service
# Install: sudo cp deploy/goes-nc-cache-cleanup.service deploy/goes-nc-cache-cleanup.timer /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable --now goes-nc-cache-cleanup.timer
#
# Runs scripts/clear_nc_cache.py, deleting raw GOES netCDF files in
# cache/goes_nc older than 1 day — they're re-downloaded on demand by
# ensure_downloaded() (app/services/goes.py) so this is pure disk-space
# reclamation, not data loss. Triggered by goes-nc-cache-cleanup.timer, not
# run directly. Shows up in Cockpit's Services page like any other systemd
# unit (search "goes-nc-cache-cleanup").
#
# Matches the convention of the existing noaa-recon-api.service,
# recon-met-update.service, and storm-archive-update.service.
[Unit]
Description=Nightly cleanup of GOES netCDF cache files older than 1 day

[Service]
Type=oneshot
WorkingDirectory=/var/www/html/noaa-recon-api
ExecStart=/var/www/html/noaa-recon-api/.venv/bin/python3 scripts/clear_nc_cache.py
User=server
Group=server

[Install]
WantedBy=multi-user.target
