# deploy/noaa-recon-api.service
# Install: sudo cp deploy/noaa-recon-api.service /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable --now noaa-recon-api.service
#
# Matches the convention of the existing recon_harvester.timer/.service on
# this host (see CLAUDE.md in the hurricanes repo).
[Unit]
Description=noaa-recon-api (GOES/TDR data API)
After=network.target

[Service]
Type=simple
User=server
WorkingDirectory=/var/www/html/noaa-recon-api
ExecStart=/var/www/html/noaa-recon-api/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --root-path /api
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
