# deploy/storm-archive-update.service
# Install: sudo cp deploy/storm-archive-update.service deploy/storm-archive-update.timer /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable --now storm-archive-update.timer
#
# Runs scripts/ingest_storms.py (HURDAT2 + ATCF b-deck ingestion — see
# app/services/storms.py for the pipeline). Triggered by
# storm-archive-update.timer, not run directly. Shows up in Cockpit's
# Services page like any other systemd unit (search "storm-archive-update").
#
# Matches the convention of the existing noaa-recon-api.service and the
# hurricanes project's recon_harvester.service/.timer on this host.
[Unit]
Description=Nightly NOAA storm track archive update (HURDAT2 + ATCF b-decks)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/var/www/html/noaa-recon-api
ExecStart=/var/www/html/noaa-recon-api/.venv/bin/python3 scripts/ingest_storms.py
User=server
Group=server

[Install]
WantedBy=multi-user.target
